# agent/linux/aegisai-agent.service.in
[Unit]
Description=AegisAI Linux Agent
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/aegisai-agent
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target